# supervisor.conf
[program:argue_websocket]
command=/usr/bin/node /var/www/argue_system/ws.njs
directory=/var/www/argue_system
restart=always
user=root
